# mpd.service
[Unit]
Description=Music Player Daemon

[Service]
ExecStart=/usr/bin/mpd --no-daemon /etc/mpd.conf
Restart=restart-always
RestartSec=1
