# sshd.service
[Unit]
Description=OpenSSH Daemon
Requires=network.target
After=network.target

[Service]
Type=forking
PIDFile=/var/run/sshd.pid
ExecStart=/etc/rc.d/sshd start
ExecStop=/etc/rc.d/sshd stop
ExecReload=/etc/rc.d/sshd reload
KillMode=process
Restart=restart-always

[Install]
WantedBy=multi-user.target
Also=network.service
