# syslog-ng.service
[Unit]
Description=System Logger Daemon
Before=syslog.target

[Service]
ExecStart=/usr/sbin/syslog-ng -d
Restart=restart-always

