# syslog-ng.service
[Unit]
Description=System Logger Daemon

[Service]
ExecStart=/usr/sbin/syslog-ng -d
Restart=restart-always
