# syslog-ng.socket
[Unit]
Description=Syslog Socket

[Socket]
ListenStream=/dev/log

[Install]
WantedBy=sockets.target
