# service/acpid@.service
[Unit]
Description=Socket Activated ACPID
After=syslog.target

[Service]
ExecStart=/usr/sbin/acpid -f
StandardInput=socket
