# service/acpid@.service
[Unit]
Description=Socket Activated ACPID

[Service]
ExecStart=/usr/sbin/acpid -f
StandardInput=socket
