# service/exim@.service
[Unit]
Description=Exim Mail Daemon per-connection server
After=syslog.target

[Service]
ExecStart=-/usr/sbin/exim -bs
StandardInput=socket
