# service/kpasswd.service
[Unit]
Description=Heimdal kpasswd Daemon

[Service]
ExecStart=/usr/sbin/kpasswdd
Restart=restart-always

[Install]
WantedBy=multi-user.target
