# service/krb5-kpropd.service
[Unit]
Description=Kerberos 5 propagation server
After=syslog.target

[Service]
ExecStart=/usr/sbin/kpropd -S

[Install]
WantedBy=multi-user.target
