# service/krb5-kdc.service
[Unit]
Description=Kerberos 5 KDC
After=syslog.target

[Service]
ExecStart=/usr/sbin/krb5kdc -n
Restart=always

[Install]
WantedBy=multi-user.target
