# service/mdadm.service
[Unit]
Description=MDADM Event Monitor
After=syslog.target

[Service]
ExecStart=/sbin/mdadm --monitor --scan

[Install]
WantedBy=multi-user.target
