# service/miredo.service
[Unit]
Description=Teredo IPv6 Tunneling Daemon

[Service]
ExecStart=/usr/sbin/miredo -f

