# service/slapd.service
[Unit]
Description=OpenLDAP server daemon

[Service]
Type=forking
EnvironmentFile=/etc/conf.d/slapd
ExecStart=/usr/sbin/slapd $SLAPD_OPTIONS

[Install]
WantedBy=multi-user.target
