# service/sshdgenkeys.service
[Unit]
Description=SSH Key Generation

[Service]
ExecStart=/usr/bin/ssh-keygen -A
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
