# service/syslog-ng.service
[Unit]
Description=System Logger Daemon

[Service]
Sockets=syslog.socket
ExecStart=/usr/sbin/syslog-ng -F
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=syslog.target
