# service/vsftpd.service
[Unit]
Description=vsftpd daemon
After=syslog.target

[Service]
ExecStart=/usr/sbin/vsftpd -olisten_ipv6=yes
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process

[Install]
WantedBy=multi-user.target
